# it-compliance-pack.conf
{
    "queries": {
        "FIREWALL-STATUS": {
             "query": "SELECT CASE WHEN COUNT(*) > 0 THEN 'Firewall Enabled' ELSE 'Firewall Disabled' END AS firewall_status FROM systemd_units WHERE id IN ('firewalld.service', 'ufw.service', 'nftables.service') AND active_state = 'active';",
            "interval": 180,
            "description": "firewall is in active state !",
            "platform": "linux",
            "version": "2.0.0",
            "snapshot": true
        },
        "ROOT-SSH-PASSWORD": {
            "query": "SELECT CASE WHEN (SELECT value FROM augeas WHERE path = '/etc/ssh/sshd_config' AND label = 'PermitRootLogin') = 'yes' AND (SELECT value FROM augeas WHERE path = '/etc/ssh/sshd_config' AND label = 'PasswordAuthentication') = 'yes' THEN 'SSH Root Password Login Enabled' ELSE 'SSH Root Password Login Disabled' END AS ssh_root_status, (SELECT CONCAT(label, ' ', value) FROM augeas WHERE path = '/etc/ssh/sshd_config' AND label = 'PermitRootLogin') AS permit_root_login, (SELECT CONCAT(label, ' ', value) FROM augeas WHERE path = '/etc/ssh/sshd_config' AND label = 'PasswordAuthentication') AS password_authentication;",
            "interval": 180,
            "description": "ssh root password is activated ! security violation .",
            "platform": "linux",
            "version": "2.0.0",
            "snapshot": true
        } 
}
}
