# it-compliance-pack-V2.conf
{
    "queries": {
        "FIREWALL-STATUS": {
             "query": "SELECT CASE WHEN COUNT(*) > 0 THEN 'Firewall Enabled' ELSE 'Firewall Disabled' END AS firewall_status FROM systemd_units WHERE id IN ('firewalld.service', 'ufw.service', 'nftables.service') AND active_state = 'active';",
            "interval": 180,
            "description": "firewall is in active state !",
            "platform": "linux",
            "version": "2.0.0",
            "snapshot": true
        },
        "ROOT-SSH-PASSWORD": {
            "query": "SELECT CASE  WHEN (SELECT COUNT(*) FROM augeas WHERE path = '/files/etc/ssh/sshd_config/PermitRootLogin' AND value = 'yes') > 0 AND (SELECT COUNT(*) FROM augeas WHERE path = '/files/etc/ssh/sshd_config/PasswordAuthentication' AND value = 'yes') > 0   THEN 'SSH Root Password Login Enabled'  ELSE 'SSH Root Password Login Disabled'  END AS ssh_root_status;",
            "interval": 180,
            "description": "ssh root password is activated ! security violation .",
            "platform": "linux",
            "version": "2.0.0",
            "snapshot": true
        }
}
